# test/systemd/modules/systemd-modules-load-sh.service
[Unit]
Description=Load Kernel Modules

[Service]
Type=oneshot
ExecStart=/etc/init.d/modutils.sh

[Install]
WantedBy=sysinit.target

