# debug/sa/stress_test/stress_test.service
[Unit]
Description=Stress Test

[Service]
ExecStart=/data/stress_test.sh

[Install]
WantedBy=multi-user.target
