# debug/sa/stress_test/stress_test_quick.service
[Unit]
Description=Stress Test Service

[Service]
Type=oneshot
ExecStart=/data/stress_test.sh

[Install]
WantedBy=sysinit.target
