# debug/sa/stress_test/target/stress_test.service
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Stress Test Service
After=multi-user.target

[Service]
ExecStartPre=/bin/sleep 5
ExecStart=/data/stress_test.sh

#ExecStartPost=/bin/sh -c "/bin/sh /data/stress_test.sh"

#Restart=always
#RestartSec=3

[Install]
WantedBy=multi-user.target
